# claude-web-console.service
[Unit]
Description=Claude Web Console
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env npx @jacexh/claude-web-console
Restart=on-failure
RestartSec=5

# Environment
Environment=PORT=3000
Environment=HOST=0.0.0.0
# Uncomment and set if claude is not in PATH
# Environment=CLAUDE_PATH=/usr/local/bin/claude
# Uncomment to set a default working directory for new sessions
# Environment=CC_WEB_CONSOLE_CWD=/home/your-user/projects

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.claude

[Install]
WantedBy=multi-user.target
